# deploy/backup/box-backup.timer
[Unit]
Description=Run box-backup nightly

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
